# supervisord.conf
[supervisord]
nodaemon=true

[program:conpot]
command=/usr/local/bin/conpot --template default --logfile /data/conpot/log/conpot.log
autorestart=true
